# nooi.service
[Unit]
Description=NOOI.net - Next.js Production Server
After=network.target

[Service]
Type=simple
User=hadmin
WorkingDirectory=/home/hadmin/nooi.net
Environment=NODE_ENV=production
Environment=PATH=/home/hadmin/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/hadmin/.local/bin/node /home/hadmin/nooi.net/node_modules/.bin/next start -p 3000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nooi

[Install]
WantedBy=multi-user.target
